# files/system/common/etc/profile.d/modern-unix.sh
# Aliases for modernized versions of common Linux tools

# bat for cat
alias cat='bat --style=plain --pager=never'

# eza for ls
alias ls='eza'
alias ll='eza -lHbg --group-directories-first'
alias la='eza -lHbgaa --group-directories-first'
